# dot_config/systemd/user/distrobox-upgrade.service
[Unit]
Description=Distrobox Automatic Upgrade

[Service]
Type=oneshot
ExecStart=/usr/bin/distrobox upgrade -a