# dot_config/systemd/user/opencode-rebuild.timer
[Unit]
Description=Weekly rebuild of opencode safe container image

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
